# http-light.service
[Unit]
Description=HTTP light service

[Service]
Type=simple
ExecStart=sh -c "~/.cargo/bin/pwm_http"

[Install]
WantedBy=default.target
